# systemd/scboot-reload.service
[Unit]
Description=React to scboot configuration change

[Service]
Type=oneshot
ExecStart=@DST_SCBOOT_SYNC_DKMS@
